# launchd/com.zhaowenlong.self-wiki-weekly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.zhaowenlong.self-wiki-weekly</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/zsh</string>
		<string>/Users/zhaowenlong/workspace/dev.self-wiki/scripts/launchd-weekly.sh</string>
	</array>
	<key>WorkingDirectory</key>
	<string>/Users/zhaowenlong/workspace/dev.self-wiki</string>
	<key>StartCalendarInterval</key>
	<dict>
		<key>Weekday</key>
		<integer>0</integer>
		<key>Hour</key>
		<integer>4</integer>
		<key>Minute</key>
		<integer>0</integer>
	</dict>
	<key>StandardOutPath</key>
	<string>/Users/zhaowenlong/workspace/dev.self-wiki/self-wiki/log/launchd-weekly.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/zhaowenlong/workspace/dev.self-wiki/self-wiki/log/launchd-weekly.err.log</string>
	<key>RunAtLoad</key>
	<false/>
</dict>
</plist>
